# deploy/systemd/bot.service
[Unit]
Description=trad-bot funding-rate arbitrage
After=network-online.target chronyd.service
Wants=network-online.target chronyd.service
# Stop restart-looping after 5 crashes in 5 min — a crash loop on a
# live trading host must page a human, not keep re-firing orders.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/trad-bot
EnvironmentFile=/opt/trad-bot/.env
ExecStart=/opt/trad-bot/.venv/bin/python -m src.main --config /opt/trad-bot/config/live.yaml
Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trad-bot

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/trad-bot/data /var/lib/bot
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
