# deploy/systemd/regime-bot.service
[Unit]
Description=trad-bot regime-switching long/short perp bot
After=network-online.target chronyd.service
Wants=network-online.target chronyd.service

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/trad-bot
EnvironmentFile=/opt/trad-bot/.env
ExecStart=/opt/trad-bot/.venv/bin/python -m src.strategy.regime_live --config /opt/trad-bot/config/regime_switch.yaml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trad-bot-regime

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/trad-bot/data /var/lib/bot
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
